# deploy/cron/lineup-sync.cron
# Poll ESPN/API-Football for official lineups before kickoff.
# Install with: deploy/scripts/install-cron.sh
#
*/5 * * * * cd /opt/world-cup-2026 && docker compose -f docker-compose.prod.yml exec -T api flask --app wsgi sync-lineups >> /var/log/wc26-lineup-sync.log 2>&1
